# openvpn-3proxy/filesystem/etc/service/openvpn/run
#!/bin/sh
echo "Starting openvpn..."
cd /etc/openvpn
exec openvpn --config /etc/openvpn/openvpn.conf 2>&1
