# openvpn-3proxy/filesystem/etc/service/openvpn/run
#!/bin/sh
echo "Starting openvpn..."

mkdir -p /dev/net
if [ ! -c /dev/net/tun ]; then
    mknod /dev/net/tun c 10 200
fi

cd /etc/openvpn
exec openvpn --config /etc/openvpn/openvpn.conf 2>&1
