# openvpn-3proxy/filesystem/etc/service/openvpn/run
#!/bin/sh
echo "Starting OpenVPN..."
exec openvpn --config /etc/openvpn/openvpn.conf 2>&1
