# openvpn-3proxy/filesystem/etc/service/3proxy/run
#!/bin/sh
echo "Starting 3proxy..."
exec /usr/local/bin/3proxy /etc/3proxy.cfg 2>&1
